# root/etc/s6-overlay/s6-rc.d/init-codex-config/run
#!/usr/bin/with-contenv bash

if [[ ! -f /usr/src/app/docs-config.yaml ]]; then
    cp /defaults/docs-config.yaml /usr/src/app/docs-config.yaml
fi

chown -R abc:abc /usr/src/app
